# config/includes.chroot/etc/skel/.config/autostart/.conkyboot.sh.desktop
[Desktop Entry]
Type=Application
Exec=/home/nat/.conkyboot.sh
Hidden=false
X-MATE-Autostart-enabled=true
Name[fr_FR]=Conky
Name=Conky
Comment[fr_FR]=Conky
Comment=Conky
