# config/includes.chroot/etc/skel/.config/autostart/conky.desktop
[Desktop Entry]
Name=conky
Exec=/bin/conky
Hidden=false
Type=Application
X-MATE-Autostart-enabled=true
