# config/includes.chroot/etc/update-motd.d/01-message
#!/bin/sh

HOSTNAME='uname -n'

echo "Bienvenue sur NATOS\n\n"
echo "Hotsname : $HOSTNAME\n\n"
echo "Uptime : $(uptime)\n\n"

